# co.codewizards.cloudstore.server/src/etc/init/cloudstore-server.conf
# cloudstore - file synchronisation server
#
# CloudStore is a file synchronisation tool and this init-conf starts its server so that
# remote clients can connect to it. See: http://cloudstore.codewizards.co

description	"cloudstore file synchronisation server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /opt/cloudstore/bin/cloudstore-server
